# .config/systemd/user/ollama-serve.service
[Unit]
Description=ollama serve

[Service]
ExecStart=/usr/bin/ollama serve
Restart=always

[Install]
WantedBy=default.target
